# custom_components/unas_pro/scripts/fan_control.service
[Unit]
Description=Fan Control Based on Temperatures
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/root/fan_control.sh --service
Restart=always
RestartSec=30
User=root

[Install]
WantedBy=multi-user.target
