# temp_monitor.service
[Unit]
Description=Temperature Monitor with MQTT
After=network.target

[Service]
ExecStart=/root/temp_monitor.sh
Restart=always
User=root
WorkingDirectory=/root

[Install]
WantedBy=multi-user.target